# scripts/systemd-template.socket
[Unit]
# DYNAMICALLY GENERATED FILE! Edit at your own risk
# compile:noescape
# This socket file supports both single instance and multi-instance deployments
# Single instance: hytale-server.socket
# Multi-instance: hytale-server@<instance-id>.socket (using systemd template)
BindsTo=%n.service

[Socket]
ListenFIFO=/var/run/%n.socket
Service=%n.service
RemoveOnStop=true
SocketMode=0660
SocketUser=$GAME_USER
User=$GAME_USER
Group=$GAME_USER
